# deploy/supervisor-nexo.conf
; ══════════════════════════════════════════════════════════════════════════════
; Nexo Digital Store — Supervisor config para Laravel Horizon
; Ruta destino: /etc/supervisor/conf.d/nexo-horizon.conf
;
; Activar:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start nexo-horizon:*
; ══════════════════════════════════════════════════════════════════════════════

[program:nexo-horizon]
process_name=%(program_name)s
command=php /var/www/nexokeys/artisan horizon
autostart=true
autorestart=true
user=nexo
redirect_stderr=true
stdout_logfile=/var/www/nexokeys/storage/logs/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
stopsignal=SIGTERM
killasgroup=true
priority=999
